# config/includes.chroot/etc/update-motd.d/00-fydelis
#!/bin/bash

echo "
███████╗██╗   ██╗██████╗ ███████╗██╗     ██╗███████╗
██╔════╝██║   ██║██╔══██╗██╔════╝██║     ██║██╔════╝
█████╗  ██║   ██║██████╔╝█████╗  ██║     ██║███████╗
██╔══╝  ██║   ██║██╔═══╝ ██╔══╝  ██║     ██║╚════██║
██║     ╚██████╔╝██║     ███████╗███████╗██║███████║
╚═╝      ╚═════╝ ╚═╝     ╚══════╝╚══════╝╚═╝╚══════╝

FydelisTechOS - Red Team Edition
"
